# gpsd-tcp.service
[Unit]
Description=gpsd-tcp

[Service]
Type=simple
ExecStart=/usr/bin/gpsd-tcp

[Install]
WantedBy=multi-user.target
